# deliveryhub.service
[Unit]
Description=DeliveryHub SaaS Node.js Server
Documentation=https://github.com/killsis/deliveryhub
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/deliveryhub/server
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=deliveryhub
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
